# openwrt-build/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#
# write gps 
echo "$(uci get system.@system[0].latitude),$(uci get system.@system[0].longitude)">/var/run/latlon.txt

exit 0
