# openwrt-addons/etc/init.d/count_reboots
#!/bin/sh /etc/rc.common

START=11

boot()
{
	local file='/www/REBOOT_COUNTER'
	local i

        if [ -e "$file" ]; then
                read i <"$file"
                i=$(( i + 1 ))
                echo $i >"$file"
        else
                echo 1 >"$file"
        fi
}
