# openwrt-addons/etc/init.d/count_reboots
#!/bin/sh /etc/rc.common

START=11

boot()
{
	. /tmp/loader
	_watch counter '/www/REBOOT_COUNTER' increment 1

	{
		echo '### [START] loaded kmodules'
		cat '/proc/modules'
		echo '### [READY] loaded kmodules'
	} >'/dev/console'
}
