# openwrt-addons/etc/init.d/count_reboots
#!/bin/sh /etc/rc.common

START=11

boot()
{
	{
		echo '### [START] loaded kmodules'
		cat '/proc/modules'
		echo '### [READY] loaded kmodules'
	} >'/dev/console'
}
