# openwrt-addons/etc/init.d/cron.user
#!/bin/sh /etc/rc.common

START=51

fresh_flashed()
{
	# defined in openwrt-build/build.sh:
	# '* * * * * /etc/init.d/cron.user boot'
	grep -sq 'cron.user boot' '/etc/crontabs/root' && return 0

	if [ -e '/tmp/loader' ]; then
		return 1
	else
		return 0
	fi
}

boot()
{
	mkdir '/tmp/running_cron.user' || return 0

	if fresh_flashed; then
		rm '/etc/crontabs/root'
		/etc/init.d/generate_kalua_loader boot
		/etc/init.d/S51crond_fff+ restart
	else
		/etc/init.d/S51crond_fff+ start
	fi

	rm -fR '/tmp/running_cron.user'
}
