# openwrt-addons/etc/init.d/end_boot
#!/bin/sh /etc/rc.common

START=99

boot()
{
	read -r _ <'/sys/block/zram0/num_writes'	# suck in error-message in dmesg

	. /tmp/loader
	. /usr/sbin/cron.check_klog
	rm '/tmp/BOOT' '/tmp/DEBUG'

	_http robotstxt_generate

	# this removes interface 'tunl0@NONE'
	rmmod ipip

	{
		echo 'cron.upgrade_packages'
		echo "rm '$TMPDIR/scheduler_keep_running'"

		local reason
		read -r reason <"$TMPDIR/REBOOT_REASON" && {
			[ "$reason" = 'nightly_reboot' ] || {
				echo '_log remote ready_booting'
			}
		}
	} >>$SCHEDULER
}
