# openwrt-addons/etc/init.d/end_boot
#!/bin/sh /etc/rc.common

START=99

boot()
{
        rm /tmp/BOOT

        scheduler -a "cron.upgrade_packages; rm /tmp/DEBUG"

	# fixme! why does crond sometimes not start? race with cron.watchdog?
	pidof crond >/dev/null || /etc/init.d/S51crond_fff+ start
}
