# openwrt-addons/etc/init.d/end_boot
#!/bin/sh /etc/rc.common

START=99

boot()
{
	. /tmp/loader
	rm '/tmp/BOOT' '/tmp/DEBUG'

	_http robotstxt_generate
	echo >>$SCHEDULER 'cron.upgrade_packages'
	echo >>$SCHEDULER "rm '$TMPDIR/scheduler_keep_running'"

	# FIXME! only work after a real crash?
#	bool_true 'system.@system[0].avoid_autoreboot' && {
#		_log it first_run daemon sms 'startup'
#	}
}
