# openwrt-addons/etc/init.d/end_boot
#!/bin/sh /etc/rc.common

START=99

boot()
{
	. /tmp/loader
	rm '/tmp/BOOT' '/tmp/DEBUG'

	_http robotstxt_generate

	# this removes interface 'tunl0@NONE'
	rmmod ipip

	echo >>$SCHEDULER 'cron.upgrade_packages'
	echo >>$SCHEDULER "rm '$TMPDIR/scheduler_keep_running'"
	echo >>$SCHEDULER '_log remote ready_booting'
}
