# openwrt-addons/etc/init.d/generate_kalua_loader
#!/bin/sh /etc/rc.common

START=00

boot()
{
	# /proc/sys/vm/panic_on_oom = 1
	# /proc/sys/kernel/panic_on_oops = 1
	# /proc/sys/kernel/panic = 10

	for ENTRY in "vm.panic_on_oom=1" "kernel.panic_on_oops=1" "kernel.panic=10"; do {
		/sbin/sysctl -w $ENTRY
		grep -q ^"$ENTRY" "/etc/sysctl.conf" || {
			echo "$ENTRY" >>"/etc/sysctl.conf"
		}
	} done

    touch /tmp/BOOT
    touch /tmp/DEBUG
    /etc/kalua_init

}

