# openwrt-addons/etc/init.d/generate_netparams
#!/bin/sh /etc/rc.common

START=41

boot()
{
	. /tmp/loader
	_netparam check
	/etc/init.d/S43ethers_fff+ start
}
