# openwrt-addons/etc/init.d/generate_kalua_loader
#!/bin/sh /etc/rc.common

START=00

boot()
{
	# /proc/sys/vm/panic_on_oom = 1
	# /proc/sys/kernel/panic_on_oops = 1
	# /proc/sys/kernel/panic = 10

	for ENTRY in "vm.panic_on_oom=1" "kernel.panic_on_oops=1" "kernel.panic=10"; do {
		/sbin/sysctl -w $ENTRY
		grep -q ^"$ENTRY" "/etc/sysctl.conf" || {
			echo "$ENTRY" >>"/etc/sysctl.conf"
		}
	} done

	/etc/init.d/zram start

        touch /tmp/BOOT
        touch /tmp/DEBUG
        /etc/kalua_init
	. /tmp/loader

	[ $( _system ram_size ) -gt 16384 ] || {
		[ -e "/etc/init.d/apply_profile" ] && {
			killall logger		# fixme! maybe patch '/etc/init.d/rcS'

			/etc/init.d/telnet disable
			/etc/init.d/dropbear disable
			/etc/init.d/dnsmasq disable

			case "$CONFIG_PROFILE" in
				*adhoc)
				;;
				*)
					uci set wireless.radio0.disabled=1
				;;
			esac
		}
	}

	[ -e "/etc/config/wireless" ] || return 0

	# patch olsrd for using likely an adhoc-interfaces if we are in hybrid mode + fallback to wlan0
	grep -q "Cell:" /etc/init.d/olsrd || {
		sed -i 's#if network_get_device IFNAME "$interface"; then$#if network_get_device IFNAME "$interface"; then case $IFNAME in wlan*) for IFNAME in $IFNAME $IFNAME-1 $IFNAME-2 $IFNAME-3 $IFNAME; do iw dev $IFNAME info | fgrep -q IBSS \&\& break; done;; esac#' /etc/init.d/olsrd
	}

	case "$CONFIG_PROFILE" in
		ffweimar*)
		;;
		schoeneck_adhoc)
			uci set wireless.radio0.htmode="40-"
		;;
		hotello-K80_ap|boltenhagen_ap|dhsylt_ap)
			uci set wireless.radio0.channel="$( _math random_integer 1 11 )"
		;;
		*)
			[ -e "/lib/modules/$( uname -r )/b43.ko" ] && {
				local file="/lib/wifi/mac80211.sh"
				local keyword="keyspec}"	# must be at the end of a line
				local command1='config_get bitrates "$device" bitrates'
				local command2='test -n "$bitrates" \&\& iw dev "$ifname" set bitrates legacy-2.4 $bitrates'

				[ "$( uci get wireless.@wifi-iface[0].mode )" = "adhoc" ] && {
					[ -n "$( uci get wireless.radio0.bitrates )" ] || {
						uci set wireless.radio0.bitrates="6 9 12 18 24 36 48 54"

						case "$( uci get wireless.@wifi-iface[0].mcast_rate )" in
							1000|2000|5500|11000)
								uci delete wireless.@wifi-iface[0].mcast_rate
							;;
						esac
					}

					grep -q "$keyword"$ "$file" && {
						sed -i "s/$keyword$/$keyword ; $command1 ; $command2 /" "$file"
					}
				}
			}
		;;
	esac

	case "$HARDWARE" in
		"Buffalo WHR-HP-G54")
			case "$( uci get wireless.radio0.rxantenna )-$( uci get wireless.radio0.txantenna )" in
				"1-1")
				;;
				*)
					uci set wireless.radio0.rxantenna=1
					uci set wireless.radio0.txantenna=1
					uci commit wireless
				;;
			esac
		;;
		"Linksys WRT54G"*)
			case "$( uci get wireless.radio0.rxantenna )-$( uci get wireless.radio0.txantenna )" in
				"0-0")
				;;
				*)
					uci set wireless.radio0.rxantenna=0
					uci set wireless.radio0.txantenna=0
					uci commit wireless
				;;
			esac
		;;
	esac
}
