# openwrt-addons/etc/init.d/generate_netparams
#!/bin/sh /etc/rc.common

START=41

boot()
{
        /etc/init.d/S41build_static_netparam start
        /etc/init.d/S43ethers_fff+ start
}
