# openwrt-addons/etc/init.d/generate_kalua_loader
#!/bin/sh /etc/rc.common

START=00

boot()
{
        touch /tmp/BOOT
        touch /tmp/DEBUG
        /etc/kalua_init
}
