# openwrt-addons/etc/init.d/generate_kalua_loader
#!/bin/sh /etc/rc.common

START=00

boot()
{
	# http://www.kernel.org/doc/Documentation/sysctl/kernel.txt
	# http://www.kernel.org/doc/Documentation/sysctl/vm.txt
	# /proc/sys/vm/panic_on_oom = 2
	# /proc/sys/kernel/panic_on_oops = 1
	# /proc/sys/kernel/panic = 10

	for ENTRY in "vm.panic_on_oom=2" "kernel.panic_on_oops=1" "kernel.panic=10"; do {
		/sbin/sysctl -w $ENTRY
		grep -q ^"$ENTRY" "/etc/sysctl.conf" || {
			echo "$ENTRY" >>"/etc/sysctl.conf"
		}
	} done

	# zram is included since r35033
	grep -q ^"zram " "/proc/modules" || /etc/init.d/zram start

	# fixme! why does it hang? with plain OpenWrt this does not happen
	fgrep -sq "killall logger" "/etc/init.d/rcS" || echo "killall logger" >>"/etc/init.d/rcS"

        touch /tmp/BOOT
        touch /tmp/DEBUG
        /etc/kalua_init
	. /tmp/loader

	# fixme!
	echo >"/tmp/dhcp.leases" "1 11:22:33:44:55:66 9.8.7.6 check_if_cron_is_running *"

	[ -n "$LOWMEM" ] && {
		[ -e "/etc/init.d/apply_profile" ] && {
			/etc/init.d/telnet disable
			/etc/init.d/dropbear disable
			/etc/init.d/dnsmasq disable
			uci set network.wan.proto=static	# dont start udhcpc
			uci set dhcp.@dnsmasq[0].disabled=1

			case "$CONFIG_PROFILE" in
				*adhoc)
				;;
				*)
					uci set wireless.radio0.disabled=1
				;;
			esac
		}
	}

	[ -e "/etc/config/wireless" ] || return 0

	# patch olsrd for using likely an adhoc-interfaces if we are in hybrid mode + fallback to wlan0
	grep -q "Cell:" /etc/init.d/olsrd || {
		sed -i 's#if network_get_device IFNAME "$interface"; then$#if network_get_device IFNAME "$interface"; then case $IFNAME in wlan*) for IFNAME in $IFNAME $IFNAME-1 $IFNAME-2 $IFNAME-3 $IFNAME; do iw dev $IFNAME info | fgrep -q IBSS \&\& break; done;; esac#' /etc/init.d/olsrd
	}
}
