# openwrt-addons/etc/init.d/netfilter
#!/bin/sh /etc/rc.common

START=45

. /tmp/loader

start() {
	_netfilter start
}

stop() {
	_netfilter stop
}
