# openwrt-addons/etc/init.d/override_uci_vars
#!/bin/sh /etc/rc.common
. /tmp/loader

START=01

boot()
{
	case "$CONFIG_PROFILE" in
		*_ap)
			# typically the mesh is on channel 11
			uci set wireless.radio0.channel="$( _math random_integer 1 7 )"
		;;
	esac
}
