# openwrt-addons/etc/init.d/prepare_stuff
#!/bin/sh /etc/rc.common

START=70

boot()
{
	/etc/init.d/S70prepare_fff+ start
	/etc/init.d/S72filter_fff+ start
}
