# openwrt-addons/etc/init.d/S99vpn
#!/bin/sh /etc/rc.common

START=99

boot() 
{

. /tmp/loader 

if [ -e /www/SOFTWARE_FULLY_INSTALLED ] && [ $(_net local_inet_offer) ]; then 
   _log do S99vpn notice "looks like we have internet.. let's start the VPN early"
   _vpn start
else 
   _log do S99vpn err "something went wrong.. no vpn at this time"
fi

} 
