# openwrt-addons/etc/init.d/set_etc_hardware
#!/bin/sh /etc/rc.common 
# set /etc/HARDWARE from machine field of /proc/cpuinfo

START=15 

boot()
{
    HARDWARE=$(grep ^machine /proc/cpuinfo | sed 's/.*: \(.*\)/\1/')

    if [ -z "$HARDWARE" ]; then
        logger -t set_etc_hardware -p user.err "error: found no machine entry in /proc/cpuinfo"
        exit 1
    else
        echo "$HARDWARE" > /etc/HARDWARE
        if [ $? -eq 0 ]; then        
           logger -t set_etc_hardware -p user.notice "wrote $HARDWARE to /etc/HARDWARE"
        else 
           logger -t set_etc_hardware -p user.err "error writing $HARDWARE to /etc/HARDWARE"
           exit 1
        fi
    fi 
}

