# openwrt-addons/etc/hotplug.d/iface/05-netfilter-wan
#!/bin/sh

[ "$ACTION" = 'ifup' -a "$INTERFACE" = 'wan' ] && {
	. /tmp/loader
	_log do hotplug_wan daemon info 'iface wan is up - starting masquerading'
        _netfilter start
}
