# openwrt-addons/etc/hotplug.d/iface/05-netfilter-wan
#!/bin/sh

. /tmp/loader

[ "$ACTION" = ifup -a "$INTERFACE" = "wan" ] && {
	_log do "iface wan is up - starting masquerading"
        _netfilter start 
}
