# openwrt-addons/etc/hotplug.d/iface/60-dyndns
#!/bin/sh

case "$ACTION" in
	ifup)
		uci -q get network.$INTERFACE.dyndns >/dev/null && {
			. /tmp/loader
			echo >>$SCHEDULER_IMPORTANT "_wget do \"$( uci -q get network.$INTERFACE.dyndns )\""
		}
	;;
esac
