# openwrt-addons/etc/hotplug.d/iface/60-dyndns
#!/bin/sh

case "$ACTION" in
	ifup)
		uci -q get network.$INTERFACE.dyndns >/dev/null && {
			. /tmp/loader
			_wget do "$( uci -q get network.$INTERFACE.dyndns )"
		}

		[ "$( uci -q get network.$INTERFACE.proto )" = 'dhcp' ] && {
			[ -z "$LODEV" ] && . /tmp/loader
			echo >>$SCHEDULER '/etc/init.d/S41build_static_netparam'
		}

		case "$INTERFACE" in
			'wan')
				[ -z "$LODEV" ] && . /tmp/loader
				WANADR="$( ubus call network.interface.wan status | grep '"address": ' | head -n1 | cut -d'"' -f4 )"
				_net ip4_is_private "$WANADR" || {
					_log do dyndns.$ACTION daemon info "prevent DNS-AmpAttac on dev '$DEVICE' for IP $WANADR"
					iptables -I INPUT -i $DEVICE -p tcp --dport 53 -j DROP
					iptables -I INPUT -i $DEVICE -p udp --dport 53 -j DROP
				}
			;;
		esac
	;;
	ifupdate|ifdown)
		# wrong place but works:
		[ -e '/etc/crontabs/root' ] && {
			. /tmp/loader
			[ "$( _stopwatch stop '/tmp/CRON_LASTCALL' interim,nolog,seconds )" -gt 300 ] && {
				# locking maybe unneeded
				mkdir '/tmp/LOCK_dhcpscript_cron_restart' && {
					if pidof crond >/dev/null; then
						_stopwatch start '/tmp/CRON_LASTCALL' global
						_log do iface_$ACTION daemon alert "cron hangs, restarting"
						/etc/init.d/S51crond_fff+ restart
					else
						_stopwatch start '/tmp/CRON_LASTCALL' global
						_log do iface_$ACTION daemon alert "cron died, starting"
						/etc/init.d/S51crond_fff+ start
					fi

					rm -fR '/tmp/LOCK_dhcpscript_cron_restart'
				}
			}
		}
	;;
esac
