# openwrt-addons/etc/hotplug.d/iface/70-qos
#!/bin/sh

case "$ACTION" in
	'ifup')
		. /tmp/loader

		interface_needs_it()
		{
			bool_true 'system.@system[0].disable_qos' && return 1

			test "$DEVICE" = "$LANDEV" && return 0
			test "$DEVICE" = "$WANDEV" && return 0
			_net dev_is_tuntap "$DEVICE" && return 0

			# TODO: maybe we do not need AP-interfaces, just adhoc/mesh
			_net dev_is_wifi "$DEVICE"
		}

		interface_needs_it && {
			if bool_true "network.$INTERFACE.shaping" ; then
				UPLINK="$(   uci -q get network.$INTERFACE.shaping_uplink )"
				DOWNLINK="$( uci -q get network.$INTERFACE.shaping_downlink )"
				_netfilter wondershaper start "$DEVICE" "$DOWNLINK" "$UPLINK"
			else
				_netfilter traffic_control_add "$DEVICE"
			fi
		}
	;;
esac
