# openwrt-addons/etc/hotplug.d/block/00-mount_storage
#!/bin/sh

case "$ACTION" in
	'add')
		case "$DEVPATH" in
			*'/block/sd'[a-z])
				. /tmp/loader
				_watch storage
			;;
		esac
	;;
esac
