# openwrt-addons/etc/hotplug.d/block/00-mount_storage
#!/bin/sh

case "$ACTION" in
	'add')
		case "$DEVPATH" in
			# /devices/platform/101c0000.ehci/usb1/1-1/1-1:1.0/host0/target0:0:0/0:0:0:0/block/sda
			# /devices/platform/101c0000.ehci/usb1/1-1/1-1:1.0/host0/target0:0:0/0:0:0:0/block/sda/sda1
			*'/block/sd'[a-z]'/sd'[a-z][0-9])
				# TODO: the first hotplug-calls are:
				# usb -> tty -> block at 8...9 secs uptime
				# and loader-generation starts at 15 sec uptime
				[ -e '/tmp/loader' ] || cp /www/loader.backup /tmp/loader

				. /tmp/loader
				_storage automount
			;;
		esac
	;;
esac
