# openwrt-addons/etc/uci-defaults/lowmem
#!/bin/sh
. /tmp/loader

[ -n "$LOWMEM" -a $OPENWRT_REV -gt 43735 ] && {
	uci set wireless.radio0.disabled='1'
	uci delete network.wan6
	uci set network.wan.proto='static'

	/etc/init.d/odhcpd disable
	/usr/sbin/uhttpd disable
	/etc/init.d/sysntpd disable
	/etc/init.d/dropbear disable
	/etc/init.d/dnsmasq disable
	/etc/init.d/log disable

	grep 'askconsole' '/etc/inittab' | grep -q ^'#' || {
		sed -i 's/.*askconsole/# &/' /etc/inittab
	}
}

false
